# etc/systemd/system/btrbk.service
[Unit]
Description=btrbk backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/btr_backup
