# etc/systemd/system/rclone.timer
[Unit]
Description=rclone daily sync

[Timer]
OnCalendar=daily
AccuracySec=10min
Persistent=true

[Install]
WantedBy=timers.target
